# solar_template.service
[Unit]
Description=My Python Script Service
After=network.target

[Service]
ExecStart=/home/pi/RPI-Solar-monitoring/start_script.sh
WorkingDirectory=/home/pi/RPI-Solar-monitoring
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=20
User=pi

[Install]
WantedBy=multi-user.target